# supervisord.conf
[supervisord]
nodaemon = true

[program:gunicorn]
command = gunicorn -k sync 'create_application:create_application("config/deploy_conf.yml")' --log-config config/logging_conf.ini --timeout 60
autostart = True
autorestart = True
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command = /usr/sbin/nginx
autostart = True
autorestart = True
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0